# blog.conf
[program:blog]
command=/usr/local/bin/gunicorn wsgi --bind 0.0.0.0:2001 --pid /tmp/blog.pid
directory=/var/www/blog
autostart=true
autorestart=true